# build/openwisp_openvpn/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/supervisord.pid

[program:openvpn]
user=root
directory=/
command=/usr/sbin/openvpn --config %(ENV_VPN_NAME)s.conf
autostart=true
autorestart=true
stopsignal=INT
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
# Set logfile maxbytes to 0 to
# avoid invalid seek error
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
