# images/openwisp_openvpn/openvpn.crontab
*/1 * * * * sh /openvpn.sh
*/5 * * * * sh /revokelist.sh
