# build/openwisp_websocket/daphne.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:daphne]
user=root
directory=/opt/openwisp/
command=/usr/local/bin/daphne -b 0.0.0.0 -p %(ENV_CONTAINER_PORT)s openwisp.asgi:channel_layer --access-log -
autostart=true
autorestart=true
stopsignal=INT
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
