# build/openwisp_websocket/worker.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:runworker]
user=root
directory=/opt/openwisp/
command=/usr/local/bin/python manage.py runworker
autostart=true
autorestart=true
stopsignal=INT
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
